# mongos/supervisord.conf
[supervisord]
nodaemon=true

[program:mongos]
command=/start_mongos.sh
autorestart=false
priority=30
