# replicaset/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/start_db.sh
autorestart=true
priority=30

[program:setup_replicaset]
command=/setup_replica.sh
autorestart=false
startsecs=0
startretries=3
